# infra/retainr-api.service
[Unit]
Description=retainr API server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=retainr
WorkingDirectory=/opt/retainr
ExecStart=/opt/retainr/api
Restart=always
RestartSec=5

# Environment
EnvironmentFile=/opt/retainr/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=retainr-api

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
